# files/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontabFrtEt8 installed on Thu Oct 27 12:30:58 2016)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
#Ansible: chown rights
* * * * * chown -R postgres /etc/postgresql/; chown -R postgres /etc/pgbouncer; chown postgres /etc/repmgr.conf > /var/log/pg_crons.log
